# templates/etc/systemd/system/data-container.service.erb
[Unit]
Description=Docker Data Container for <%= @title %>
Requires=docker.service<% if @depends %><% @depends.each do |dep| %><%= " #{dep}" %><% end %><% end %>
After=docker.service<% if @depends %><% @depends.each do |dep| %><%= " #{dep}" %><% end %><% end %>

[Service]
Type=oneshot
Restart=no
RemainAfterExit=yes
ExecStart=-/usr/bin/docker create \
    <%= @docker_run_options %> \
    <%= @image %>
ExecStop=/usr/bin/docker stop <%= @title %>

[Install]
WantedBy=multi-user.target
