# templates/etc/systemd/system/exec-container.service.erb
[Unit]
Description=Docker Exec Container for <%= @title %>
Requires=docker.service
After=docker.service <%= depends %>

[Service]
Type=oneshot
Restart=no
ExecStart=/usr/bin/docker exec -i <%= @container %> <%= @command %>

[Install]
WantedBy=multi-user.target
