# templates/etc/systemd/system/exec-container.service.erb
[Unit]
Description=Docker Exec for <%= @title %>
Requires=docker.service <%= @depends %>
After=docker.service <%= @depends %>

[Service]
Type=oneshot
Restart=on-failure
RestartSec=5
RemainAfterExit=yes
<% if @systemd_env_file %><%= "EnvironmentFile=#{@systemd_env_file}" %><% end %>
ExecStart=/usr/bin/docker exec -i <%= @container %> <%= @command %>

[Install]
WantedBy=multi-user.target
