# templates/etc/systemd/system/run-container.service.erb
[Unit]
Description=Docker Container Service for <%= @title %>
Requires=docker.service<% if @depends %><% @depends.each do |dep| %><%= " docker-#{dep}.service" %><% end %><% end %>
After=docker.service<% if @depends %><% @depends.each do |dep| %><%= " docker-#{dep}.service" %><% end %><% end %>

[Service]
Type=simple
Restart=always
RestartSec=5
EnvironmentFile=/etc/docker-<%= @title %>.env
ExecStartPre=-/usr/bin/docker stop <%= @title %>
ExecStartPre=-/usr/bin/docker rm <%= @title %>
ExecStart=/usr/bin/docker run --rm \
    <%= @docker_run_options %> \
    $IMAGE_ARG $COMMAND
ExecStop=/usr/bin/docker stop <%= @title %>

[Install]
WantedBy=multi-user.target
